# systemd/calvinbtc-watchdog.service
[Unit]
Description=CalvinBTC · Watchdog — System Supervisor & Auto-Restart (Binance)
After=network.target redis.service calvinbtc-signal.service calvinbtc-executor.service
Wants=redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/calvinbot-binance
EnvironmentFile=/opt/calvinbot-binance/.env
ExecStart=/opt/calvinbot-binance/.venv/bin/python watchdog_binance.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
